# task2/etc/supervisord.d/httpserver.ini
[program:httpserver]
command=/app/httpserver
directory=/app
startsecs=0
stopwaitsecs=60
stopsignal=QUIT
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

